# com.trellis.garden.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Copy to ~/Library/LaunchAgents/com.trellis.garden.plist and replace
     /path/to/trellis with the absolute path to your trellis checkout
     (launchd does not expand ~ or environment variables). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trellis.garden</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/path/to/trellis/run-garden.sh</string>
    </array>
    <!-- Run nightly at 03:00. If the Mac is asleep, launchd runs it at next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/path/to/trellis/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/trellis/launchd.err.log</string>
</dict>
</plist>
